# fs/etc/profile.d/more-modern-aliases.sh
# shellcheck disable=SC2148
alias .h='cd $(zoxide query --interactive)'
alias cat='bat --style plain'
alias ls='lsd'
alias lt='lsd -l --tree'
