# fs/etc/profile.d/postgres-socket.sh
# shellcheck disable=SC2148
export PGHOST=/run/postgresql
